# config/systemd/user/copyq.service
[Unit]
Description=Copyq Daemon

[Service]
Type=simple
ExecStart=copyq
ExecStop=copyq exit
Restart=always

[Install]
WantedBy=default.target

